# config/contabo-sm-139/systemd/authelia.service
[Unit]
Description=Authelia SSO (аутентификация)
After=network.target traefik.service

[Service]
Type=simple
User=authelia
Group=authelia
ExecStart=/usr/local/bin/authelia --config /etc/authelia/configuration.yml
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/authelia /var/lib/authelia /var/log/authelia
PrivateTmp=true

[Install]
WantedBy=multi-user.target
