# config/contabo-sm-139/systemd/n8n.service.example
[Unit]
Description=n8n workflow automation
Documentation=https://docs.n8n.io
After=network.target

[Service]
Type=simple
User=n8n
Environment="NODE_ENV=production"
Environment="N8N_BASIC_AUTH_ACTIVE=true"
Environment="N8N_BASIC_AUTH_USER=admin"
Environment="N8N_BASIC_AUTH_PASSWORD={{N8N_AUTH_PASSWORD}}"
Environment="N8N_HOST=127.0.0.1"
Environment="N8N_PORT=5678"
Environment="N8N_PROTOCOL=http"
Environment="WEBHOOK_URL=https://n8n.dev.borisovai.ru/"
Environment="N8N_METRICS=true"
Environment="N8N_LOG_LEVEL=info"
Environment="N8N_USER_FOLDER=/var/lib/n8n"
ExecStart=/usr/bin/n8n start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n

[Install]
WantedBy=multi-user.target
